# build/Lede_source_R3P/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
# 
if [[ `grep -c "coremark" /etc/crontabs/root` -eq 1 ]]; then
  sed -i /coremark/d /etc/crontabs/root
fi
exit 0
