# build/Lede_source_lxc/files/etc/crontabs/root
50 6 * * 1 sleep 5 && reboot
0 5 * * 6 lua /usr/share/passwall/rule_update.lua log > /dev/null 2>&1 &
0 6 * * * /usr/share/adbyby/adblock.sh > /tmp/adupdate.log 2>&1
